# ansible/roles/livegrep-index/templates/livegrep-refresh-index.service.j2
[Unit]
Description=Reload livegrep index: {{livegrep_index}}

[Service]
ExecStart=/home/livegrep/refresh-index {{livegrep_index}}
User=livegrep
Restart=no
