# ansible/roles/livegrep-indexer/templates/livegrep-reindex.service.j2
[Unit]
Description=Reindex Livegrep

[Service]
ExecStart=/home/livegrep/systemd-reindex
User=livegrep
Restart=no
